# subscription.service
[Unit]
Description=VPN Subscription API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/YadrenoVPN-main
ExecStart=/usr/bin/python3 /root/YadrenoVPN-main/subscription_api.py
Restart=always
RestartSec=10

# Логирование
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
